# nodeone.service
[Unit]
Description=Easy NodeOne (Gunicorn → wsgi:app)
After=network.target

[Service]
Type=simple
User=nodeone
Group=nodeone
WorkingDirectory=/var/www/nodeone
EnvironmentFile=-/var/www/nodeone/.env
Environment="PATH=/var/www/nodeone/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PORT=8000"
ExecStartPre=/var/www/nodeone/venv/bin/python3 /var/www/nodeone/backend/bootstrap_nodeone.py
ExecStart=/var/www/nodeone/venv/bin/gunicorn -w 3 -b 127.0.0.1:8000 --timeout 120 --access-logfile - --error-logfile - wsgi:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
