# backend/yappy-verification.service
[Unit]
Description=Verificación automática de pagos Yappy cada 5 minutos
After=network.target

[Service]
Type=oneshot
ExecStart=/home/relaticpanama2025/projects/membresia-relatic/venv/bin/python3 /home/relaticpanama2025/projects/membresia-relatic/backend/verify_yappy_payments_cron.py
StandardOutput=append:/home/relaticpanama2025/projects/membresia-relatic/logs/yappy_verification.log
StandardError=append:/home/relaticpanama2025/projects/membresia-relatic/logs/yappy_verification.log
User=relaticpanama2025
WorkingDirectory=/home/relaticpanama2025/projects/membresia-relatic/backend

[Install]
WantedBy=multi-user.target
