# backend/yappy-verification.service
[Unit]
Description=Verificación automática de pagos Yappy cada 5 minutos
After=network.target

[Service]
Type=oneshot
ExecStart=/var/www/nodeone/venv/bin/python3 /var/www/nodeone/backend/verify_yappy_payments_cron.py
StandardOutput=append:/var/www/nodeone/logs/yappy_verification.log
StandardError=append:/var/www/nodeone/logs/yappy_verification.log
User=nodeone
WorkingDirectory=/var/www/nodeone/backend

[Install]
WantedBy=multi-user.target
